# cron_files/pi_at_b29.crontab
@reboot	sleep 60; $HOME/projects/arduino-interaction/calculateLengthOfNight/startPICommunications.py /dev/ screen_name
