# cron_files/pi_at_b52testpi.crontab
@reboot	sleep 60; $HOME/projects/arduino-interaction/calculateLengthOfNight/startPICommunications.py /dev/ttyACM0 arduino_screen_1
@reboot	sleep 70; $HOME/projects/arduino-interaction/calculateLengthOfNight/startPICommunications.py /dev/ttyACM1 arduino_screen_2
