# cron_files/pi_at_buzon.crontab
*/5 * * * * ~/dynudns/dynu.sh >/dev/null
