# cron_files/pi_at_alorberpi1.crontab
@reboot	$HOME/bin/cleanup_chromium_preferences
#@hourly $HOME/bin/cleanup_chromium_preferences