# cron_files/pi_at_b29testpi.crontab
@reboot	sleep 60; $HOME/projects/arduino-interaction/calculateLengthOfNight/startPICommunications.py /dev/ttyUSB0 arduino_screen
