# cron_files/root_at_turbo.crontab
@reboot	sleep 60; /Users/aal/sagwa_bin/loop_openvpn_on_turbo
